# .config/systemd/user/swww-daemon.service
[Unit]
Description=swww wallpaper daemon

[Service]
ExecStart=/usr/bin/swww-daemon
Restart=always

[Install]
WantedBy=default.target

