# .config/systemd/user/wallpaper-sort.service
[Unit]
Description=Sort wallpapers into day/night folders

[Service]
Type=oneshot
ExecStart=/home/sarang/scripts/sort-wallpapers.sh
